# fan-control.service
[Unit]
Description=Control the fan

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-control

[Install]
WantedBy=default.target

